# terraform/ansible/files/zramswap.service
[Unit]
Description=Enable ZRAM

[Service]
ExecStart=/bin/bash /usr/sbin/zramswap.sh

[Install]
WantedBy=multi-user.target
